# web21.conf
[program:web21]
command=/usr/local/bin/gunicorn wsgi --bind localhost:2000 --workers 3 --worker-class gevent
directory=/var/www/web21
autostart=true
autorestart=true

stdout_logfile=/var/log/supervisor/web21_stdout.log
stderr_logfile=/var/log/supervisor/web21_stderr.log